# c.conf
[program:c]
command=python3 /home/curator2/c/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/c/test.err.log
stdout_logfile=/home/curator2/c/test.out.log